# dev.pydantic.mactab.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.pydantic.mactab</string>
    <!-- __BIN__ is replaced with the installed executable path by `make startup-install`. -->
    <key>ProgramArguments</key>
    <array>
        <string>__BIN__</string>
    </array>
    <!-- Launch at login. Not KeepAlive: if you Quit from the menu bar it
         stays quit until the next login, which is the expected behaviour. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
